# wp6-virtualfolder/rc.local
mkdir -p /home/vagrant/logs
nohup mono /home/vagrant/build/MetadataService/MetadataService.exe > /home/vagrant/logs/MetadataService.log 2>&1 &
nohup /home/vagrant/scripts/startVRE.sh > /home/vagrant/logs/VRE.log 2>&1 &
