# wp6-virtualfolder/bootstrapcloud/rc.local
mkdir -p /home/vagrant/logs
/home/vagrant/scripts/startMetaDataService.sh
/home/vagrant/scripts/startVRE.sh
/home/vagrant/scripts/startVirtuoso.sh
