# wp6-virtualfolder/bootstrap/rc.local
mkdir -p /home/vagrant/logs
/home/vagrant/scripts/startMetaDataService.sh
