# wp6-virtualfolder/bootstrap/rc.local
mkdir -p /home/vagrant/logs
sudo -H -u vagrant /home/vagrant/scripts/startMetaDataService.sh
/home/vagrant/scripts/startVRE.sh
